# deploy/systemd/gladiator-api.service
[Unit]
Description=Gladiator API (Fastify + Socket.IO)
After=network.target

[Service]
Type=simple
WorkingDirectory=/opt/gladiator/apps/server
Environment=PORT=5174
Environment=CORS_ORIGIN=https://web.example.com
Environment=AP_REGEN_MS=6000
ExecStart=/usr/bin/npm run dev
Restart=always
RestartSec=2

[Install]
WantedBy=multi-user.target


